# bin/videoloop.cron
PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin

*/5 * * * * (cd /var/pivideoloop/pi-video-loop && bin/environment node index.js >> /var/pivideoloop/videoloop.log)
